# runner/launchd/com.muse.brain.orchestrator.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.muse.brain.orchestrator</string>

    <key>ProgramArguments</key>
    <array>
      <string>__RUNNER_DIR__/run-orchestrator.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__RUNNER_DIR__</string>

    <key>StartInterval</key>
    <integer>60</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>__LOG_DIR__/orchestrator.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>__LOG_DIR__/orchestrator.stderr.log</string>
  </dict>
</plist>
